# rc.local
#!/bin/sh
/mnt/wdc/router/ddnspod.sh
/www/server/redis6378/src/redis-server  /www/server/redis6378/redis.conf
exit 0

~

